# containers/files/systemd/unpack-skills.service
[Unit]
Description=Unpacks skills mounted as OCI artifacts

[Service]
Type=oneshot
ExecStart=/usr/local/bin/unpack-skills.sh
